# switcher/systemd/vohive-web-gateway.service
[Unit]
Description=VoHive web gateway with USB mode controls
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/vohive-web-gateway
Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
